# project1/etcd.service.j2
[Unit]
Description=etcd

[Service]
ExecStart=/tmp/etcd-download-test/etcd
Restart=always

[Install]
WantedBy=multi-user.target
